# config/autostart/xrandr.desktop
[Desktop Entry]
Encoding=UTF-8
Type=Application
Name=XRANDR
Comment=
Exec=xrandr --fb 1600x900
StartupNotify=false
Terminal=false
Hidden=false
